# tdmf_ose/unix_v1.2/ftdsrc/agent/startagent.sh.Linux
#!/bin/sh
# LICENSED MATERIALS / PROPERTY OF IBM
# %PRODUCTNAME% version %VERSION%
# (c) Copyright %COMPANYNAME%  %COPYRIGHTYEAR 2001%.  All Rights Reserved.
# The source code for this program is not published or otherwise
# divested of its trade secrets, irrespective of what has been
# deposited with the U.S. Copyright Office.
# US Government Users Restricted Rights - Use, duplication or
# disclosure restricted by GSA ADP Schedule Contract with %COMPANYNAME%
#
#
# This script will start the agent daemons for the %PRODUCTNAME% 
# product.
#
# Generally, this should run in init mode 2-5 (in the /etc/rc[2-5].d directory),
# but can be moved, if needed in a local installation, to other places after
# the network has started.  It must start after the other %PKGNM%-startagent script.
#
# exit values: 0 (OK), 1 (FAIL), and 2 (N/A)

# %PKGNM% agent Script

# chkconfig: 2345 26 26
# description: %PKGNM% in.%QAGN% Script
#

# /etc/init.d/%PKGNM%-startagent

### BEGIN INIT INFO
# Provides:       %PKGNM%-startagent
# Required-Start: $network $syslog
# Should-Start: $portmap %PKGNM%-startdaemons
# Required-Stop:
# Should-Stop:  $network $syslog
%ifdistribution == redhat
# Default-Start:  2 3 4 5
%endif
%ifdistribution == suse
# Default-Start:  2 3 5
%endif
# Default-Stop:
# Description:    $COMPANYNAME2% %PRODUCTNAME% Agent Script
### END INIT INFO
#
LANG=C
export LANG
PATH="/sbin:/bin:/usr/sbin:/usr/bin:$PATH"
export PATH


case "$1" in
'start')		
    /%OPTDIR%/%PKGNM%/bin/launchagent
    result=$?
%ifdistribution == redhat
    if [ $result -eq 0 ]
    then
        if [ -d /var/lock/subsys ] ; then
	    touch /var/lock/subsys/%PKGNM%-startagent
        fi
    fi
%endif
    ;;
'stop')
    /%OPTDIR%/%PKGNM%/bin/killagent
%ifdistribution == redhat
    if [ -e /var/lock/subsys/%PKGNM%-startagent ] ; then
	rm /var/lock/subsys/%PKGNM%-startagent
    fi
%endif
    ;;
*)
    echo "Usage: $0 { start | stop }"
    exit 1
    ;;
esac

exit 0
